# deploy/systemd/bot.service.example
[Unit]
Description=CitaConsulares Bot (Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/CitaConsulares/bot
Environment=NODE_ENV=production
EnvironmentFile=%h/CitaConsulares/bot/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

