# deploy/systemd/worker.service.example
[Unit]
Description=CitaConsulares Worker (Playwright)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/CitaConsulares/worker
Environment=NODE_ENV=production
EnvironmentFile=%h/CitaConsulares/worker/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

